# ops/launchd/at.dieter.hq.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>at.dieter.hq</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/dieter/.openclaw/workspace/dieter-hq</string>

    <!-- Use a login shell so Node/npm from your normal PATH (nvm/volta/homebrew) works. -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>./scripts/hq-prod.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOST</key>
      <string>127.0.0.1</string>
      <key>PORT</key>
      <string>3010</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/dieter-hq.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dieter-hq.launchd.err.log</string>
  </dict>
</plist>
